# systemd/ogma-dream.timer
[Unit]
Description=Run the Ogma memory consolidation nightly at ~02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
